# .config/metapkg/jckimble-laptop.install
post_install(){
	post_upgrade
}

post_upgrade(){
	holo apply
	systemctl enable iptables
	systemctl enable wol-disable
	systemctl enable NetworkManager
}
